# wsl2-ssh-agent.service.in
[Unit]
Description=Share Windows ssh-agent to WSL2 using SSH agent forwarding
After=network-online.target

[Service]
Type=idle
Environment=RELAY_MIN_COUNT=0 SSHD_PORT=34672
ExecStart=<exec>
RestartSec=3
TimeoutStopSec=1
Restart=on-failure

[Install]
WantedBy=default.target
